# supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:jira-api]
command=python -m jira_emulator serve --host 0.0.0.0 --port %(ENV_PORT)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true

[program:mcp-server]
command=python /app/mcp_servers/atlassian_jira.py
environment=JIRA_SERVER="http://localhost:%(ENV_PORT)s",JIRA_USER="%(ENV_JIRA_USER)s",JIRA_TOKEN="%(ENV_JIRA_TOKEN)s",MCP_PORT="%(ENV_MCP_PORT)s",MCP_HOST="0.0.0.0"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
